# umnico-greeting-bot.service
[Unit]
Description=Umnico Greeting Bot (Webhook)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/umnico-greeting-bot
ExecStart=/usr/bin/python3 /opt/umnico-greeting-bot/app_webhook.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/umnico-greeting-bot/bot.log
StandardError=append:/opt/umnico-greeting-bot/bot.log

[Install]
WantedBy=multi-user.target
